# universal-mkvtoolnix/root/etc/s6-overlay/s6-rc.d/init-mod-universal-mkvtoolnix-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-sabnzbd-mkvtoolnix-add-package/run
